# queries/windows/scheduled_tasks_creation.yml
title: Scheduled Tasks Creation
description: Detection of schtasks /create command as well as any use of the New-ScheduledTask* 
  cmdlets from powershell.
author: keyboardcrunch && kryssar
date: 2020-10-10
modified: 2025-05-07
mitre:
   tactic: Execution, Persistence
   technique: T1053
   subtechnique: 005
operating_system: windows
query: ((tgt.process.name = "shtasks.exe" AND tgt.process.cmdline contains "/create") OR (src.process.cmdline contains "New-ScheduledTask" OR cmdScript.content contains "New-ScheduledTask")) AND !(src.process.parent.name in ("services.exe", "OfficeClickToRun.exe"))
false_positives: 
   - NPCAP driver task
   - Legitimate task creation
tags: 

